# examples/incursion_update.timer
[Unit]
Description=Run incursion updater every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
